# hades/others/mouse/mouse.rules
SUBSYSTEM=="input", ACTION=="add", ATTR{name}=="G603 Mouse", RUN+="/home/andre.rodrigues/.local/bin/mouse.sh"